# scripts/launchd/ai.diting.prefetch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.diting.prefetch</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/<dev-user>/diting-radar/.venv/bin/python</string>
    <string>/Users/<dev-user>/diting-radar/scripts/prefetch-vault.py</string>
  </array>
  <key>StartCalendarInterval</key>
  <array>
    <dict><key>Hour</key><integer>9</integer><key>Minute</key><integer>50</integer></dict>
    <dict><key>Hour</key><integer>13</integer><key>Minute</key><integer>50</integer></dict>
    <dict><key>Hour</key><integer>17</integer><key>Minute</key><integer>50</integer></dict>
    <dict><key>Hour</key><integer>19</integer><key>Minute</key><integer>50</integer></dict>
  </array>
  <key>StandardOutPath</key>
  <string>/Users/<dev-user>/diting-radar/state/cron-prefetch.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/<dev-user>/diting-radar/state/cron-prefetch.log</string>
</dict>
</plist>
